# deploy/systemd/ocr-api.service
[Unit]
Description=FastAPI backend for PDFTools OCR
After=network.target

[Service]
WorkingDirectory=/opt/pdftools/backend
Environment=PYTHONPATH=/opt/pdftools/backend
EnvironmentFile=/opt/pdftools/.env
ExecStart=/opt/pdftools/venv/bin/uvicorn app.main:app --host ${API_HOST} --port ${API_PORT}
Restart=always

StandardOutput=append:/opt/pdftools/backend/logs/ocr.log
StandardError=append:/opt/pdftools/backend/logs/ocr.log

[Install]
WantedBy=multi-user.target
