# install/systemd/ocr-api.service
[Unit]
Description=FastAPI backend for PDFTools OCR
After=network.target

[Service]
WorkingDirectory=/opt/pdftools/backend
Environment=PYTHONPATH=/opt/pdftools/backend
ExecStart=/usr/bin/python3 -m uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always

StandardOutput=append:/opt/pdftools/backend/logs/ocr.log
StandardError=append:/opt/pdftools/backend/logs/ocr.log

[Install]
WantedBy=multi-user.target
