# install/systemd/purge-ocr.timer
[Unit]
Description=Planification horaire de la purge OCR

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
 