# install/systemd/purge-ocr.service
[Unit]
Description=Purge des jobs OCR expirés

[Service]
Type=oneshot
WorkingDirectory=/opt/pdftools/backend
Environment=PYTHONPATH=/opt/pdftools/backend
ExecStart=/usr/bin/python3 scripts/purge_old_jobs.py
